# sf/systemd/awecron.service
[Unit]
Description=Awecron cron daemon

[Service]
ExecStart=/bin/sh /usr/local/bin/awecron
User=root

[Install]
WantedBy=multi-user.target
